# 2025082616/script/queue/queue_processor.conf
[program:doye_queue_processor]
; Supervisor配置文件 - DoyeCMS AI队列处理器
; 用于确保队列处理器长期稳定运行

; 程序基本配置
command=php d:\phpstudy_pro\WWW\DoyeCMS\script\queue\queue_processor.php
directory=d:\phpstudy_pro\WWW\DoyeCMS\script\queue
user=www
autostart=true
autorestart=true
startsecs=10
startretries=3

; 日志配置
stdout_logfile=d:\phpstudy_pro\WWW\DoyeCMS\logs\supervisor_queue_stdout.log
stderr_logfile=d:\phpstudy_pro\WWW\DoyeCMS\logs\supervisor_queue_stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_capture_maxbytes=1MB

; 进程管理
numprocs=1
process_name=%(program_name)s_%(process_num)02d
priority=999

; 环境变量
environment=PATH="/usr/local/bin:/usr/bin:/bin",HOME="/tmp"

; 停止配置
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true

; 重启策略
exitcodes=0,2
restart_when_binary_changed=true

; 其他配置
redirect_stderr=false
stdout_events_enabled=false
stderr_events_enabled=false